# roles/alloy/templates/alloy.service.d/override.conf.j2
[Service]
{% if alloy_docker_enabled | bool %}
# The embedded cAdvisor exporter needs container runtime socket access that the
# package's unprivileged alloy user does not have on Docker/containerd hosts.
User=root
Group=root
{% else %}
User=alloy
Group=alloy
{% endif %}
{% if alloy_supplemental_groups | length > 0 %}
SupplementaryGroups={{ alloy_supplemental_groups | join(' ') }}
{% else %}
SupplementaryGroups=
{% endif %}
