# roles/certbot/templates/systemd/certbot-renew.service.j2
[Unit]
Description=Renew Certbot certificates

[Service]
Type=oneshot
ExecStart={{ certbot_renew_command }}
